# nextcloud-valkey.container
[Container]
AutoUpdate=registry
Pod=cloud.pod
Image=docker.io/valkey/valkey

[Service]
TimeoutStartSec=900

